# systemd/actual-ai-categorizer.service
# One-shot run, triggered by the companion .timer.
# Installed by scripts/install.sh. Single-directory layout: code, config,
# secrets, logs and data all live under /opt/actual-ai-categorizer.
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now actual-ai-categorizer.timer
[Unit]
Description=Actual Budget AI auto-categorizer (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=actual
WorkingDirectory=/opt/actual-ai-categorizer
# Secrets referenced as ${VAR} in config.yaml are loaded from this file:
EnvironmentFile=/opt/actual-ai-categorizer/secrets.env
# Paths in config.yaml are relative to WorkingDirectory above.
ExecStart=/usr/bin/node dist/index.js run --config config.yaml
